# deployment/banking-frontend.service
[Unit]
Description=Banking Chatbot Frontend (Streamlit)
After=network.target banking-backend.service
Wants=banking-backend.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/banking-chatbot
EnvironmentFile=/home/ubuntu/banking-chatbot/backend/.env
Environment=PYTHONUNBUFFERED=1
Environment=API_URL=http://127.0.0.1:8000
Environment=STREAMLIT_SERVER_HEADLESS=true
Environment=STREAMLIT_BROWSER_GATHER_USAGE_STATS=false

ExecStart=/home/ubuntu/banking-chatbot/venv/bin/streamlit run \
    /home/ubuntu/banking-chatbot/frontend/streamlit_app.py \
    --server.address 127.0.0.1 \
    --server.port 8501 \
    --server.enableCORS false \
    --server.enableXsrfProtection true

Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/banking-chatbot/logs/frontend.log
StandardError=append:/home/ubuntu/banking-chatbot/logs/frontend.log

MemoryMax=300M

[Install]
WantedBy=multi-user.target
